# templates/service-systemd.service.j2
#jinja2: lstrip_blocks: True
# {{ ansible_managed }}

[Unit]
Description={{service.description|default('Autogenerated service for ' + service.name)}}

[Service]
{% if service.exec_start_pre is defined %}
ExecStartPre={{service.exec_start_pre}}
{% endif %}
ExecStart={{service.exec_start}}
Restart={{service.restart|default('always')}}
RestartSec={{service.restart_sec|default(2)}}
{% if service.user is defined %}
User={{service.user}}
{% endif %}
{% for env in service.environment|default({}) %}
Environment={{env}}={{service.environment[env]}}
{% endfor %}


[Install]
WantedBy={{service.wanted_by|default('multi-user.target')}}